# wg-monitor@.service
[Unit]
Description=WireGuard Monitor for %i
Documentation=https://github.com/chouyu/wg_monitor
After=network-online.target wg-quick@%i.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# 加载实例特定的配置文件
EnvironmentFile=/etc/default/wg-monitor-%i

# 工作目录
WorkingDirectory=/opt/wg-monitor

# 执行命令 (使用环境变量)
ExecStart=/bin/sh -c 'exec /usr/bin/python3 /opt/wg-monitor/wg_monitor.py \
    --log-path "${LOG_PATH}" \
    --interval "${INTERVAL}" \
    --threshold "${THRESHOLD}" \
    --stats-interval "${STATS_INTERVAL}" \
    $([ "${DEBUG}" = "true" ] && echo "--debug")'

# Reload (SIGHUP) - 用于 logrotate 后重新打开日志文件
ExecReload=/bin/kill -HUP $MAINPID

# 重启策略
Restart=always
RestartSec=10s
StartLimitInterval=300
StartLimitBurst=5

# 资源限制
LimitNOFILE=1024
LimitNPROC=10
MemoryMax=256M
CPUQuota=20%

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true

# 能力限制
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wg-monitor-%i

[Install]
WantedBy=multi-user.target
