# wg-monitor.service
[Unit]
Description=WireGuard Peer Status Monitor
Documentation=https://github.com/chouyu/wg_monitor
After=network-online.target wg-quick@wg0.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# 加载环境变量 (如果文件不存在，systemd 会报错，这是预期的，因为我们需要配置)
EnvironmentFile=/etc/default/wg-monitor

# 工作目录
WorkingDirectory=/opt/wg-monitor

# 执行命令 (使用环境变量)
ExecStart=/bin/sh -c 'exec /usr/bin/python3 /opt/wg-monitor/wg_monitor.py \
    --log-path "${LOG_PATH}" \
    --interval "${INTERVAL}" \
    --threshold "${THRESHOLD}" \
    --stats-interval "${STATS_INTERVAL}" \
    $([ "${DEBUG}" = "true" ] && echo "--debug") \
    $([ -n "${INTERFACE}" ] && echo "--interface ${INTERFACE}") \
    $([ -n "${ON_CHANGE_SCRIPT}" ] && echo "--on-change-script ${ON_CHANGE_SCRIPT}")'

# Reload (SIGHUP) - 用于 logrotate 后重新打开日志文件
ExecReload=/bin/kill -HUP $MAINPID

# 重启策略
Restart=always
RestartSec=10s
StartLimitInterval=300
StartLimitBurst=5

# 超时配置
TimeoutStartSec=30s
TimeoutStopSec=30s

# 资源限制
LimitNOFILE=1024
LimitNPROC=10
MemoryMax=256M
CPUQuota=20%

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictNamespaces=true

# 能力限制
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wg-monitor

[Install]
WantedBy=multi-user.target
